# systemd/smu-auth.service
[Unit]
Description=SMU Auth Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/smu-microservices/services/auth-service
Environment="DATABASE_URL=postgres://smuuser:smupass@localhost:5432/smu_auth"
Environment="JWT_SECRET=your-super-secret-jwt-key-change-this-in-production"
Environment="PORT=8001"
ExecStart=/opt/smu-microservices/services/auth-service/auth-service
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
