# systemd/smu-auth.service
[Unit]
Description=SMU Auth Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/smu-microservices/services/auth-service
Environment="PORT=8001"
Environment="DB_HOST=localhost"
Environment="DB_PORT=5432"
Environment="DB_USER=smuuser"
Environment="DB_PASSWORD=smupass"
Environment="DB_NAME=smu_auth"
Environment="JWT_SECRET=your-super-secret-jwt-key-change-this-in-production"
Environment="JWT_EXPIRATION=24h"
Environment="REDIS_URL=redis://localhost:6379"
Environment="RABBITMQ_URL=amqp://guest:guest@localhost:5672/"
ExecStart=/opt/smu-microservices/services/auth-service/auth-service
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
